# defshine/0024/awesome-flask-todo/etc/awesome-flask-todo.conf
[program:awesome-flask-todo]
command = /home/defshine/www/awesome-flask-todo/venv/bin/gunicorn -b 127.0.0.1:9000 wsgi_gunicorn:app
directory = /home/defshine/www/awesome-flask-todo
user = defshine
